# edison/mbm3000.service
[Unit]
Description=MBM3000 service

[Service]
Type=simple
ExecStart=/home/root/run.sh

[Install]
WantedBy=default.target
